# com.sni.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sni.pipeline</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/scott/.bun/bin/bun</string>
		<string>/Users/scott/Projects/sni-research-v2/scripts/pipeline.js</string>
		<string>--mode</string>
		<string>full</string>
	</array>

	<key>StartCalendarInterval</key>
	<array>
		<!-- Thursday -->
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/scott/Projects/sni-research-v2</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/scott/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/scott/Projects/sni-research-v2/logs/pipeline.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/scott/Projects/sni-research-v2/logs/pipeline-error.log</string>

	<key>Nice</key>
	<integer>10</integer>
</dict>
</plist>
