# com.sni.ev-extract.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sni.ev-extract</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/scott/.bun/bin/bun</string>
        <string>scripts/ev-link-extract.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/scott/Projects/sni-research-v2</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/scott/Projects/sni-research-v2/logs/ev-extract.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/scott/Projects/sni-research-v2/logs/ev-extract-error.log</string>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
